# scripts/aegiscore.service
[Unit]
Description=AegisCore — Adaptive Linux Threat Intelligence Platform
After=network-online.target auditd.service
Wants=network-online.target auditd.service

[Service]
Type=simple
# Dedicated kullanıcı: sudo useradd -r -s /bin/false siem
# Log dosyaları için: sudo usermod -aG adm siem
User=siem
Group=siem
WorkingDirectory=/opt/aegiscore
EnvironmentFile=-/etc/default/aegiscore
ExecStart=/opt/aegiscore/.venv/bin/python /opt/aegiscore/main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aegiscore
UMask=0027

# Güvenlik sertleştirme
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectProc=invisible
ProcSubset=pid
ReadWritePaths=/opt/aegiscore/data /var/log
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
MemoryDenyWriteExecute=false

[Install]
WantedBy=multi-user.target
